# wired/ead.service.in
[Unit]
Description=Ethernet service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman.ead
ExecStart=@libexecdir@/ead
LimitNPROC=1
Restart=on-failure
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectHome=yes
ProtectSystem=strict
ProtectControlGroups=yes
ProtectKernelModules=yes
ConfigurationDirectory=ead
ReadWritePaths=@wired_storagedir@

[Install]
WantedBy=multi-user.target
