# wired/ead.service.in
[Unit]
Description=Ethernet service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman.ead
ExecStart=@libexecdir@/ead
NotifyAccess=main
LimitNPROC=1
Restart=on-failure
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectHome=yes
ProtectSystem=strict
ProtectControlGroups=yes
ProtectKernelModules=yes
ConfigurationDirectory=ead
StateDirectory=ead
StateDirectoryMode=0700

[Install]
WantedBy=multi-user.target
