# src/iwd.service.in
[Unit]
Description=Wireless service
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=net.connman.iwd
ExecStart=@libexecdir@/iwd
LimitNPROC=1
Restart=on-failure

[Install]
WantedBy=multi-user.target
