# src/iwd.service.in
[Unit]
Description=Wireless service

[Service]
ExecStart=@libexecdir@/iwd
LimitNPROC=1

[Install]
WantedBy=multi-user.target
